# data/supervisord.conf
[supervisord]

pidfile=%(ENV_SUPERVISORD_RUNDIR)s/supervisord.pid
logfile=%(ENV_SUPERVISORD_LOGDIR)s/supervisord.log
childlogdir=%(ENV_SUPERVISORD_LOGDIR)s/

minfds=%(ENV_SUPERVISORD_MINFDS)s               ; number of startup file descriptors
minprocs=%(ENV_SUPERVISORD_MINPROCS)s           ; number of process descriptors

logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace

nodaemon=false                                  ; run supervisord as a daemon
user=root                                       ; default user


[program:overlord]
user=druid
command=java
  -server

  -Xmx%(DRUID_OVERLORD_Xmx)s
  -Xms%(DRUID_OVERLORD_Xms)s
  -XX:NewSize=%(DRUID_OVERLORD_XXNewSize)s
  -XX:MaxNewSize=%(DRUID_OVERLORD_XXMaxNewSize)s
  -XX:+UseConcMarkSweepGC
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Djava.io.tmpdir=/mnt/tmp
  -Ddruid.service=druid/overlord

  -Ddruid.host=%(ENV_HOST_IP)s

  -Ddruid.extensions.loadList=[\"%(DRUID_STORAGE_TYPE)-metadata-storage\"]
  -Ddruid.extensions.directory=%(ENV_DRUID_ROOT)s/extensions
  -Ddruid.extensions.hadoopDependenciesDir=%(ENV_DRUID_ROOT)s/hadoop_dependencies

  -Ddruid.metadata.storage.type=%(DRUID_STORAGE_TYPE)
  -Ddruid.metadata.storage.connector.connectURI=jdbc:%(DRUID_STORAGE_TYPE)://%(DRUID_STORAGE_CONN_STR)
  -Ddruid.metadata.storage.connector.user=%(DRUID_STORAGE_USER)
  -Ddruid.metadata.storage.connector.password=%(DRUID_STORAGE_PASS)

  -Ddruid.indexer.autoscale.doAutoscale=true
  -Ddruid.indexer.autoscale.strategy=ec2
  -Ddruid.indexer.autoscale.workerIdleTimeout=PT90m
  -Ddruid.indexer.autoscale.terminatePeriod=PT5M

  -Ddruid.indexer.logs.type=s3
  -Ddruid.indexer.logs.s3Bucket=druid
  -Ddruid.indexer.logs.s3Prefix=prod/logs/v1
  -Ddruid.indexer.runner.type=remote
  -Ddruid.indexer.storage.type=metadata
  -Ddruid.indexer.queue.startDelay=PT0M

  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server overlord
numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/overlord.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/overlord.log

[program:middlemanager]
user=druid
command=java
  -server
  -Xmx%(DRUID_MIDDLEMANAGER_Xmx)s
  -Xms%(DRUID_MIDDLEMANAGER_Xms)s
  -XX:+UseConcMarkSweepGC
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Djava.io.tmpdir=/mnt/tmp


  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.service=druid/middlemanager

  -Ddruid.indexer.logs.type=s3
  -Ddruid.indexer.logs.s3Bucket=%(DRUID_INDEXER_LOGS_S3_BUCKET)s
  -Ddruid.indexer.logs.s3Prefix=%(DRUID_INDEXER_LOGS_S3_PREFIX)s
  -Ddruid.indexer.task.baseTaskDir=%(DRUID_INDEXER_BASE_TASK_DIR)s

  -Ddruid.indexer.fork.property.druid.monitoring.monitors=["com.metamx.metrics.JvmMonitor"]
  -Ddruid.indexer.fork.property.druid.processing.buffer.sizeBytes=536870912
  -Ddruid.indexer.fork.property.druid.processing.numThreads=%(DRUID_PROCESSING_NUM_TREADS)s
  -Ddruid.indexer.fork.property.druid.segmentCache.locations="[{\"path\":\"%(DRUID_SEGMENTCACHE_LOCATION)s\",\"maxSize\":0}]"
  -Ddruid.indexer.fork.property.druid.server.http.numThreads=%(DRUID_SERVER_HTTP_NUM_TREADS)s
  -Ddruid.indexer.fork.property.druid.storage.archiveBaseKey=%(DRUID_STORAGE_ARCHIVE_BASE_KEY)s
  -Ddruid.indexer.fork.property.druid.storage.archiveBucket=aws-prod-druid-archive
  -Ddruid.indexer.fork.property.druid.storage.baseKey=%(DRUID_STORAGE_BASE_KEY)s
  -Ddruid.indexer.fork.property.druid.storage.bucket=%(DRUID_STORAGE_BUCKET)s
  -Ddruid.indexer.fork.property.druid.storage.type=s3

  -Ddruid.worker.capacity=9
  -Ddruid.worker.ip=%(ENV_HOST_IP)s

  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server middleManager

numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/middlemanager.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/middlemanager.log







[program:druid-coordinator]
user=druid
command=java
  -server
  -Xmx1g
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.extensions.loadList=[\"%(DRUID_STORAGE_TYPE)-metadata-storage\"]
  -Ddruid.extensions.directory=%(ENV_DRUID_ROOT)s/extensions
  -Ddruid.extensions.hadoopDependenciesDir=%(ENV_DRUID_ROOT)s/hadoop_dependencies
  -Ddruid.metadata.storage.type=%(DRUID_STORAGE_TYPE)
  -Ddruid.metadata.storage.connector.connectURI=jdbc:%(DRUID_STORAGE_TYPE)://%(DRUID_STORAGE_CONN_STR)
  -Ddruid.metadata.storage.connector.user=%(DRUID_STORAGE_USER)
  -Ddruid.metadata.storage.connector.password=%(DRUID_STORAGE_PASS)
  -Ddruid.coordinator.startDelay=PT5S
  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server coordinator
numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/coordinator.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/coordinator.log


[program:druid-historical]
user=druid
command=java
  -server
  -Xmx1g
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.extensions.loadList=[\"druid-s3-extensions\"]
  -Ddruid.extensions.directory=%(ENV_DRUID_ROOT)s/extensions
  -Ddruid.extensions.hadoopDependenciesDir=%(ENV_DRUID_ROOT)s/hadoop_dependencies
  -Ddruid.s3.accessKey=%(ENV_ACCESS_KEY)s
  -Ddruid.s3.secretKey=%(ENV_SECRET_KEY)s
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.segmentCache.locations="[{\"path\":\"/var/tmp/druid/indexCache\",\"maxSize\":5000000000}]"
  -Ddruid.server.maxSize=5000000000
  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server historical
numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/historical.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/historical.log

[program:druid-broker]
user=druid
command=java
  -server
  -Xmx1g
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.broker.cache.sizeInBytes=33554432
  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server broker
redirect_stderr=true
priority=100
[program:broker]
user=druid
command=java
  -server
  -Xmx%(DRUID_BROKER_Xmx)s
  -Xms%(DRUID_BROKER_Xms)s
  -XX:NewSize=%(DRUID_BROKER_XXNewSize)s
  -XX:MaxNewSize=%(DRUID_BROKER_XXMaxNewSize)s
  -XX:MaxDirectMemorySize=%(DRUID_BROKER_XXMaxDirectMemorySize)s
  -XX:+UseConcMarkSweepGC
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Djava.io.tmpdir=/mnt/tmp
  -Dcom.sun.management.jmxremote.port=17071
  -Dcom.sun.management.jmxremote.authenticate=false
  -Dcom.sun.management.jmxremote.ssl=false

  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.service=druid/broker
  -Ddruid.processing.buffer.sizeBytes=%(DRUID_PROCESSING_BUFFER_SIZEBYTES)s
  -Ddruid.processing.numThreads=%(DRUID_PROCESSING_NUM_TREADS)s
  -Ddruid.server.http.numThreads=%(DRUID_SERVER_HTTP_NUM_TREADS)s
  -Ddruid.broker.http.readTimeout=PT5M
  -Ddruid.broker.http.numConnections=20
  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server broker

numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/broker.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/broker.log

[program:realtime]
user=druid
command=java
  -server
  -Xmx%(DRUID_REALTIME_Xmx)s
  -Xms%(DRUID_REALTIME_Xms)s
  -XX:NewSize=%(DRUID_REALTIME_XXNewSize)s
  -XX:MaxNewSize=%(DRUID_REALTIME_XXMaxNewSize)s
  -XX:MaxDirectMemorySize=%(DRUID_REALTIME_XXMaxDirectMemorySize)s
  -XX:+UseConcMarkSweepGC
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -XX:+HeapDumpOnOutOfMemoryError
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
  -Djava.io.tmpdir=/mnt/tmp
  -Dcom.sun.management.jmxremote.port=17071
  -Dcom.sun.management.jmxremote.authenticate=false
  -Dcom.sun.management.jmxremote.ssl=false

  -Ddruid.host=%(ENV_HOST_IP)s
  -Ddruid.service=druid/realtime
  -Ddruid.processing.buffer.sizeBytes=%(DRUID_PROCESSING_BUFFER_SIZEBYTES)s
  -Ddruid.processing.numThreads=%(DRUID_PROCESSING_NUM_TREADS)s
  -Ddruid.server.http.numThreads=%(DRUID_SERVER_HTTP_NUM_TREADS)s
  -Ddruid.broker.http.readTimeout=PT5M
  -Ddruid.broker.http.numConnections=20
  -Ddruid.monitoring.monitors=["io.druid.segment.realtime.RealtimeMetricsMonitor", "com.metamx.metrics.JvmMonitor"]
  -cp %(ENV_DRUID_ROOT)s/lib/*
  io.druid.cli.Main server realtime

numprocs                        = 1
numprocs_start                  = 0
directory                       = %(ENV_DRUID_ROOT)s
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
redirect_stderr                 = true
stdout_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/realtime.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = %(ENV_SUPERVISORD_LOGDIR)s/realtime.log
